# app/data/devoirs/TP_Gentoo/evaluation.tar/tmp/evaluation/cron.d/backup.sh
37 13 * * * root /bin/tar czf /mnt/backup/data.tar.gz /hote/toto/data
